# data/limba-cleanup.service.in
[Unit]
Description=Remove uninstalled software from disk and free other resources
DefaultDependencies=no
Before=shutdown.target

[Service]
Type=oneshot
ConditionPathExists=/var/lib/limba/cleanup-needed
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/lipa cleanup --no-fail

[Install]
WantedBy=shutdown.target
