# data/limba.service.in
[Unit]
Description=Limba helper daemon

[Service]
Type=dbus
BusName=org.freedesktop.Limba
User=root
ExecStart=@CMAKE_INSTALL_FULL_LIBDIR@/limba/limba-daemon
